# watcher/systemd/alternance-watcher.service
[Unit]
Description=Alternance Discord Watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/alternance-discord-watcher
EnvironmentFile=%h/alternance-discord-watcher/.env
ExecStart=%h/alternance-discord-watcher/.venv/bin/python -m watcher.main --once
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
